# deploy/examples/nq-publish.service
[Unit]
Description=NQ witness publisher (serves GET /state)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=nq
Group=nq
UMask=0027
ExecStart=/usr/local/bin/nq-witness --config /etc/nq/publisher.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
