# deploy/examples/nq-serve.service
[Unit]
Description=NQ monitor aggregator and web UI
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=nq
Group=nq
UMask=0027
WorkingDirectory=/var/lib/nq
ExecStart=/usr/local/bin/nq-monitor serve --config /etc/nq/aggregator.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
